# SOURCES/renewCerts.service
[Unit]
Description=Renew Certificates from lets encrypt

[Service]
Type=oneshot
ExecStart=/bin/certbot renew
