# SOURCES/renewCerts.service
[Unit]
Description=Renew Certificates from lets encrypt
Requires=network-online.target

[Service]
Type=oneshot
ExecStart=/bin/certbot renew
